# hf_supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:backend]
command=/usr/local/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/home/user/app/backend
environment=DATABASE_URL="postgresql://user@127.0.0.1:5432/atlas_db",REDIS_URL="redis://127.0.0.1:6379",MONGO_URL="mongodb://127.0.0.1:27017",PYTHONUNBUFFERED=1
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/tmp/backend.log
stderr_logfile=/tmp/backend_error.log

[program:frontend]
command=/usr/bin/npm start
directory=/home/user/app/frontend
environment=PORT=7860,NEXT_PUBLIC_API_URL=""
autostart=true
autorestart=true
stdout_logfile=/tmp/frontend.log
stderr_logfile=/tmp/frontend_error.log