# hf_supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server --bind 127.0.0.1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log

[program:mongodb]
command=mongod --dbpath /data/db --bind_ip 127.0.0.1
user=mongodb
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mongodb.log

[program:postgres]
command=/bin/bash -c "sudo -u postgres /usr/lib/postgresql/*/bin/postgres -D /var/lib/postgresql/*/main"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log

[program:backend]
command=/usr/local/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
environment=DATABASE_URL="postgresql://atlas:atlaspass@127.0.0.1:5432/atlas_db",REDIS_URL="redis://127.0.0.1:6379",MONGO_URL="mongodb://127.0.0.1:27017",PYTHONUNBUFFERED=1
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_error.log

[program:frontend]
command=/bin/bash -c "cd /app/frontend && /usr/bin/npm start"
environment=PORT="3000",NEXT_PUBLIC_API_URL=""
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend_error.log